# etc/skel/.config/autostart/activate-bandwidth-widget.desktop
[Desktop Entry]
Name=activate-bandwidth-widget
Exec=sudo /usr/share/xevuan/activate-bandwidth-widget
Hidden=false


